# scripts/pi-assistant.service
[Unit]
Description=Pi Assistant Backend
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/pi-assistant
EnvironmentFile=%h/pi-assistant/.env
ExecStart=%h/pi-assistant/.venv/bin/python run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
